# ansible/roles/systemd/templates/home-stack.service.j2
[Unit]
Description=home-stack Docker Compose
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ deploy_user }}
Group={{ deploy_user }}
SupplementaryGroups=docker
WorkingDirectory=/home/{{ deploy_user }}/home-stack
Environment=HOME_STACK_ENV_FILE=/etc/home-stack.env
EnvironmentFile=/etc/home-stack.env
ExecStart=/usr/bin/docker compose -f docker/docker-compose.yml up -d --remove-orphans
ExecReload=/usr/bin/docker compose -f docker/docker-compose.yml up -d --wait --remove-orphans
ExecStop=/usr/bin/docker compose -f docker/docker-compose.yml down

[Install]
WantedBy=multi-user.target
